# system/wm/etc/profile.d/qt.sh
[ "$XDG_CURRENT_DESKTOP" = "KDE" ] || [ "$XDG_CURRENT_DESKTOP" = "GNOME" ] || export QT_STYLE_OVERRIDE="kvantum"
