# system/wm-general/etc/profile.d/fix-qt-styles.sh
# This makes Qt5 mimick GTK+ styles
export XDG_CURRENT_DESKTOP=GNOME
