# system/wm-general/etc/profile.d/java-nonreparenting-fix.sh
# Fix applications not receiving size updates when using a non reparenting
# window manager
export _JAVA_AWT_WM_NONREPARENTING=1
