# system/wm-general/etc/profile.d/java-nonreparenting-fix.sh
# Fix applications not receiving size updates when using a non reparenting
# window manager

# Re-enable this for bspwm
# i3 is a reparenting WM so it's not required, and it messes with HiDPI
#export _JAVA_AWT_WM_NONREPARENTING=1
