# user/xfce/.config/systemd/user/xfce4-notifyd.service
[Unit]
Description=XFCE notifications service

[Service]
Type=dbus
BusName=org.freedesktop.Notifications
ExecStart=/usr/lib/xfce4/notifyd/xfce4-notifyd
Restart=always

[Install]
WantedBy=default.target
